# config/systemd-service-templates/lagrange-onebot.service
[Unit]
Description=Lagrange.OneBot Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/lagrange
ExecStart=/opt/lagrange/Lagrange.OneBot
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lagrange-onebot

# 环境变量
Environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT=1

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lagrange

# 资源限制
MemoryMax=512M
TasksMax=100

[Install]
WantedBy=multi-user.target
