# config/systemd-service-templates/qq-bot.service
[Unit]
Description=QQ Bot Service
After=network.target lagrange-onebot.service
Wants=network-online.target
Requires=lagrange-onebot.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/botbwiki
Environment=PATH=/home/ubuntu/botbwiki/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ubuntu/botbwiki/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qq-bot

# 安全设置
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# 资源限制
MemoryMax=256M
TasksMax=50

[Install]
WantedBy=multi-user.target
